# builder/src/templates/debian/postinst
#!/bin/sh

systemctl daemon-reload || true
systemctl enable b2ghald || true
systemctl enable capyloon || true
systemctl start b2ghald || true
systemctl start capyloon || true
