# builder/src/templates/debian/prerm
#!/bin/sh

systemctl stop capyloon || true
systemctl stop b2ghald || true
systemctl disable b2ghald || true
systemctl daemon-reload || true
